# .local/profile.d/brainux.sh
#!/bin/sh

PATH=~/git/brainux/bin:$PATH
